# core/etc/profile.d/editor.sh
#!/bin/sh
export EDITOR=nvim
export VISUAL=nvim
