# arch/server/etc/systemd/system/backup-btrfs.timer
[Unit]
Description=Automated snapshotting of the samba filesystem.

[Timer]
OnCalendar=%i
Persistent=true

[Install]
WantedBy=timers.target
