# arch/server/etc/systemd/system/backup-btrfs.timer
[Unit]
Description=Automated daily snapshotting of the btrfs filesystem.

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
