# arch/server/etc/systemd/system/backup-azure.service
[Unit]
Description=rclone samba to azure cloud storage

[Service]
Type=oneshot
ExecStart=rclone sync -P /srv/samba azure-icebox-backup:/srv/samba
