# arch/server/etc/systemd/system/backup-azure.timer
[Unit]
Description=Daily rclone samba to azure cloud storage

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
