# arch/server/etc/systemd/system/ddns.timer
[Unit]
Description=Update icebox.fresh4less.org DNS hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
