# arch/server/etc/systemd/system/post-samba-updates.timer
[Unit]
Description=Post changes to the Samba filesystem to Slack

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
