# server/etc/systemd/system/pacman-checkupdates.service
[Unit]
Description=update pacman checkupdates database for motd
[Service]
Type=oneshot
ExecStart=checkupdates
