# server/etc/systemd/system/post-vaccine-updates.service
[Unit]
Description=Post changes to vaccine HTML to Slack

[Service]
Type=oneshot
ExecStart=post-vaccine-updates
