# server/etc/systemd/system/post-samba-updates.service
[Unit]
Description=Post changes to the Samba filesystem to Slack
After=backup-btrfs.service

[Service]
Type=oneshot
ExecStart=post-samba-updates
