# server/etc/systemd/system/post-vaccine-updates.timer
[Unit]
Description=Post changes to vaccine HTML to Slack

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
